# minimal-rpi-os/DEBIAN/postinst
#!/bin/bash

# Possibly run setup on first boot of the machine
#systemctl enable minimal-rpi-os-setup
