# minimal-rpi-os/etc/udev/rules.d/50-brcmfmac.rules
SUBSYSTEM=="net", DRIVERS=="brcmfmac", ACTION=="add", RUN+="/usr/sbin/iw '%E{INTERFACE}' set power_save off"
