# minimal-ec2-os/DEBIAN/postinst
#!/bin/bash

# Possibly run setup on first boot of the machine
systemctl enable minimal-ec2-os-setup

# Fix the NTP server to AWS' own EC2 NTP server
sed -i -E 's|^#?NTP=.*$|NTP=169.254.169.123|g' /etc/systemd/timesyncd.conf

# Mask/disable TTYs on console and enable serial
systemctl enable "serial-getty@ttyS0.service"
systemctl disable "getty@tty1.service"
systemctl mask "getty@tty1.service"
